# deploy/systemd/nrg-commerce-backup.timer
[Unit]
Description=Run NRG Commerce encrypted PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
Persistent=true
RandomizedDelaySec=15m
Unit=nrg-commerce-backup.service

[Install]
WantedBy=timers.target
